# deploy/app.service.example
# Example systemd unit to keep `next start` running on the Ubuntu host.
# Copy to /etc/systemd/system/nbu-application.service, edit paths/user, then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now nbu-application

[Unit]
Description=NBU online application system (Next.js)
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/nbu-application
ExecStart=/usr/bin/npm run start -- -p 3000
Restart=on-failure
RestartSec=5
EnvironmentFile=/var/www/nbu-application/.env

[Install]
WantedBy=multi-user.target
